# cgateweb/rootfs/etc/services.d/cgate-setup/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Configure C-Gate CNI interface after C-Gate starts (managed mode only)
#
# On first run: creates the project and adds the CNI network
# On subsequent runs: skips setup if project tag file already exists
# ==============================================================================

CGATE_MODE=$(bashio::config 'cgate_mode' 'remote')
CNI_IP=$(bashio::config 'cni_ip' '')

if [[ "${CGATE_MODE}" != "managed" ]] || [[ -z "${CNI_IP}" ]]; then
    if [[ "${CGATE_MODE}" == "managed" ]] && [[ -z "${CNI_IP}" ]]; then
        bashio::log.info "C-Gate managed mode: no CNI IP configured, skipping CNI setup"
    fi
    exec sleep infinity
fi

CGATE_DIR="/data/cgate"
CGATE_PORT=$(bashio::config 'cgate_port' '20023')
CGATE_PROJECT=$(bashio::config 'cgate_project' 'HOME')
CNI_NETWORK=$(bashio::config 'cni_network' '254')

# Check if the project already exists (tag file present from a previous run)
PROJECT_TAG="${CGATE_DIR}/tag/${CGATE_PROJECT}.xml"
if [[ -f "${PROJECT_TAG}" ]]; then
    bashio::log.info "Project '${CGATE_PROJECT}' already exists (${PROJECT_TAG}), skipping CNI setup"
    bashio::log.info "C-Gate will load the project via project.start in config"
    exec sleep infinity
fi

# Wait for C-Gate to accept connections
bashio::log.info "Waiting for C-Gate to accept connections on port ${CGATE_PORT}..."
for i in $(seq 1 30); do
    if nc -z 127.0.0.1 "${CGATE_PORT}" 2>/dev/null; then
        break
    fi
    sleep 2
done

if ! nc -z 127.0.0.1 "${CGATE_PORT}" 2>/dev/null; then
    bashio::log.error "C-Gate did not start within 60 seconds, cannot configure CNI"
    exec sleep infinity
fi

# Give C-Gate a moment to fully initialize after port opens
sleep 3

# First run: create project and configure CNI network
bashio::log.info "First-time setup: creating project '${CGATE_PROJECT}' with CNI ${CNI_IP} on network ${CNI_NETWORK}"

RESPONSE=$(
{
    sleep 1
    echo "project new ${CGATE_PROJECT}"
    sleep 1
    echo "project start ${CGATE_PROJECT}"
    sleep 1
    echo "project use ${CGATE_PROJECT}"
    sleep 1
    echo "net create ${CNI_NETWORK} cni ${CNI_IP}"
    sleep 1
    echo "project save ${CGATE_PROJECT}"
    sleep 1
    echo "quit"
} | nc 127.0.0.1 "${CGATE_PORT}" 2>&1
)

bashio::log.info "C-Gate CNI setup response:"
echo "${RESPONSE}" | while IFS= read -r line; do
    bashio::log.info "  ${line}"
done

bashio::log.info "CNI configuration complete — C-Gate will discover devices from the network"

exec sleep infinity
