# cgateweb/rootfs/etc/services.d/cgate/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Run C-Gate server (managed mode only)
# ==============================================================================

CGATE_MODE=$(bashio::config 'cgate_mode' 'remote')

if [[ "${CGATE_MODE}" != "managed" ]]; then
    bashio::log.info "C-Gate mode is '${CGATE_MODE}', C-Gate service sleeping"
    exec sleep infinity
fi

CGATE_DIR="/data/cgate"
CGATE_JAR="${CGATE_DIR}/cgate.jar"

if [[ ! -f "${CGATE_JAR}" ]]; then
    bashio::log.error "C-Gate not installed. Check cont-init.d logs."
    sleep 30
    exit 1
fi

CGATE_PORT=$(bashio::config 'cgate_port' '20023')
CGATE_EVENT_PORT=$(bashio::config 'cgate_event_port' '20025')
CGATE_PROJECT=$(bashio::config 'cgate_project' 'HOME')

bashio::log.info "Starting C-Gate server..."
bashio::log.info "  Command port: ${CGATE_PORT}"
bashio::log.info "  Event port: ${CGATE_EVENT_PORT}"
bashio::log.info "  Project: ${CGATE_PROJECT}"

# Configure ports and settings in C-GateConfig.txt
# C-Gate does not accept port/gui args on the command line
CGATE_CONFIG="${CGATE_DIR}/config/C-GateConfig.txt"
mkdir -p "${CGATE_DIR}/config"

# Build config ensuring required settings are present
declare -A SETTINGS=(
    ["cgate.port"]="${CGATE_PORT}"
    ["cgate.event.port"]="${CGATE_EVENT_PORT}"
    ["project.default"]="${CGATE_PROJECT}"
    ["project.start"]="${CGATE_PROJECT}"
)

for KEY in "${!SETTINGS[@]}"; do
    VALUE="${SETTINGS[$KEY]}"
    if [[ -f "${CGATE_CONFIG}" ]] && grep -q "^${KEY}=" "${CGATE_CONFIG}"; then
        sed -i "s|^${KEY}=.*|${KEY}=${VALUE}|" "${CGATE_CONFIG}"
    else
        echo "${KEY}=${VALUE}" >> "${CGATE_CONFIG}"
    fi
done

bashio::log.debug "C-GateConfig.txt contents:"
if [[ -f "${CGATE_CONFIG}" ]]; then
    while IFS= read -r line; do
        bashio::log.debug "  ${line}"
    done < "${CGATE_CONFIG}"
else
    bashio::log.warning "C-GateConfig.txt not found at ${CGATE_CONFIG}"
fi

bashio::log.debug "access.txt contents:"
ACCESS_FILE="${CGATE_DIR}/config/access.txt"
if [[ -f "${ACCESS_FILE}" ]]; then
    while IFS= read -r line; do
        bashio::log.debug "  ${line}"
    done < "${ACCESS_FILE}"
else
    bashio::log.debug "  (not found — C-Gate will use defaults)"
fi

bashio::log.debug "Existing project tag files:"
ls "${CGATE_DIR}/tag/"*.xml 2>/dev/null | while IFS= read -r f; do
    bashio::log.debug "  $(basename "${f}")"
done

cd "${CGATE_DIR}" || exit 1

bashio::log.debug "Launching C-Gate from ${CGATE_DIR} with -s flag"

exec java \
    -Djava.awt.headless=true \
    -jar "${CGATE_JAR}" \
    -s
